# docs/templates/alfred-nc-db-backup.timer
[Unit]
Description=Daily Nextcloud DB backup at 03:05

[Timer]
OnCalendar=*-*-* 03:05:00
Persistent=true

[Install]
WantedBy=timers.target
